# roles/summit/templates/supervisord_uwsgi.conf.j2
[program:summit]
command=/usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/summit.ini
autostart=true
autorestart=true
